# spin/files/etc/init.d/spin
#!/bin/sh /etc/rc.common

START=85

SPIN_SERVER_PID="/var/run/spin_mqtt.pid"

run_spin_server() {
    spin_mqtt&
    echo $! > ${SPIN_SERVER_PID}
}

apply_spin_config() {
    cd /usr/lib/spin
    /usr/lib/spin/show_ips.lua -o /etc/spin/ignore.list
    spin_config ignore load /etc/spin/ignore.list
    spin_config block load /etc/spin/block.list
    spin_config except load /etc/spin/allow.list
}

start() {
    insmod /usr/lib/spin/spin.ko
    apply_spin_config
    run_spin_server
}

stop() {
    if [ -e ${SPIN_SERVER_PID} ]; then
        kill `cat ${SPIN_SERVER_PID}`
        rm ${SPIN_SERVER_PID}
    fi
    rmmod spin
}

restart() {
    stop
    start
}

reload() {
    stop
    start
}
