# spin/files/etc/init.d/spin
#!/bin/sh /etc/rc.common

START=85

SPIN_SERVER_PID="/var/run/spin_collectd.pid"

run_spin_server() {
    cd /usr/lib/spin
    lua ./spin_collectd.lua&
    echo $! > ${SPIN_SERVER_PID}
}

start() {
    run_spin_server
}

stop() {
    if [ -e ${SPIN_SERVER_PID} ]; then
        kill `cat ${SPIN_SERVER_PID}`
        rm ${SPIN_SERVER_PID}
    fi
}

restart() {
    stop
    start
}

reload() {
    stop
    start
}

